# systemd/com.nanoclaw.service
[Unit]
Description=NanoClaw Personal Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
Environment=HOME={{HOME}}
Environment=PATH={{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart={{BUN_PATH}} {{PROJECT_ROOT}}/dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
